# roles/server-nextcloud/templates/nextcloud-watchstorage.service.j2
[Unit]
Description=Nextcloud watch storage job

[Service]
ExecStart=/usr/local/bin/server_nextcloud_watchstorage
KillMode=process
