# roles/server-nextcloud/templates/nextcloud-cron.service.j2
[Unit]
Description=Nextcloud cron.php job

[Service]
User={{ nextcloud_user }}
ExecStart=/usr/bin/{{ php_command }} -c {{ php_ini_path }} -f {{ nextcloud_path }}/cron.php
KillMode=process
