# roles/server-docker/templates/docker.service.j2
[Unit]
Description=Docker compose based {{ item.id }} server
Wants=network.target

[Service]
{% for current in item.system_depends | default([]) %}
Requires={{ current }}.service
{% endfor %}
WorkingDirectory=/opt/docker/{{ item.id }}/
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/docker compose start
ExecStop=/bin/docker compose stop

[Install]
WantedBy=multi-user.target
