# roles/server-webdav/templates/systemd/cleanup-webshares.service.j2
[Unit]
Description=Cleanup of expired files in webshare

[Service]
Type=oneshot
ExecStart=/usr/local/bin/server_webdav_cleanfiles
