# roles/pacmatic/templates/pacmatic.service.j2
[Unit]
Description=Update pacman caches and anylse changes

[Service]
Type=oneshot
ExecStart=/usr/bin/cron-pacmatic admin@{{ domain }}
