# roles/storage-management/files/backup.service
[Unit]
Description=Backup
 
[Service]
Type=oneshot
User=root
ExecStart=bash /usr/local/bin/storage_backup auto
