# roles/storage-management/files/backup.timer
[Unit]
Description=Run backup scripts nightly
 
[Timer]
OnCalendar=*-*-* 1:30:00
Persistent=true
 
[Install]
WantedBy=timers.target
