# roles/storage-management/files/storage-maintenance-mini.service
[Unit]
Description=Snapraid sync
 
[Service]
Type=oneshot
User=root
WorkingDirectory=/
ExecStart=bash /usr/local/bin/storage_maintenance full
