# roles/storage-management/files/storage-monitor.timer
[Unit]
Description=Storage monitoring of mounted pools

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
