# roles/storage-management/files/storage-backup-check.timer
[Unit]
Description=Check of all backups for timeliness

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
